# deploy/systemd/vortexops-queue.service
[Unit]
Description=VortexOps Laravel Queue Worker
After=network.target mysql.service

[Service]
Type=simple
User=__APP_USER__
Group=__APP_GROUP__
Restart=always
RestartSec=5
WorkingDirectory=__APP_DIR__
ExecStart=/usr/bin/php __APP_DIR__/artisan queue:work --sleep=3 --tries=3 --timeout=120

# Keep the worker in the same production environment as the web app.
Environment=APP_ENV=production

[Install]
WantedBy=multi-user.target
